# src/main/resources/archetype-resources/src/main/assembly/dist/bin/__artifactId__.service
[Unit]
Description=${name} Service

[Service]
ExecStart=/usr/local/${artifactId}/bin/${artifactId}-service-start

User=${artifactId}
Group=${artifactId}

[Install]
WantedBy=multi-user.target
