# src/main/resources/archetype-resources/src/main/assembly/dist/install/__artifactId__.service
[Unit]
Description=${name} Service

[Service]
ExecStart=/bin/java \
   -Dlogback.configurationFile=/etc/opt/dans.knaw.nl/${artifactId}/logback-service.xml \
   -Dapp.home=/opt/dans.knaw.nl/${artifactId} \
   -jar /opt/dans.knaw.nl/${artifactId}/bin/${artifactId}.jar run-service
# Java returns 143 even if the SIGTERM was handled correctly.
SuccessExitStatus=143

User=${artifactId}
Group=${artifactId}

[Install]
WantedBy=multi-user.target
